# conf/systemd-scrutiny-collector.timer
[Unit]
Description=Scrutiny Collector timer

[Timer]
OnCalendar=daily
OnStartupSec=60min
RandomizedDelaySec=60min
Persistent=true

[Install]
WantedBy=timers.target
